# deployment/service.service
# systemd unit for TService. Update user and paths before installation.

[Unit]
Description=TService
After=network.target

[Service]
Type=simple
User=yourusername
WorkingDirectory=/path/to/TService
ExecStart=/path/to/TService/.venv/bin/python src/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
